# deploy/btc/systemd/leiting-btc-review.service
[Unit]
Description=Leiting BTC review job
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/leiting-btc
EnvironmentFile=-/opt/leiting-btc/btc_config/live.env
ExecStart=/opt/leiting-btc/scripts/btc/run_review.sh

[Install]
WantedBy=multi-user.target
