# deploy/btc/systemd/leiting-btc-user-stream.service
[Unit]
Description=Leiting BTC Binance user stream service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/leiting-btc
EnvironmentFile=-/opt/leiting-btc/btc_config/live.env
ExecStart=/opt/leiting-btc/scripts/btc/run_user_stream_service.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
