# deploy/btc/systemd/leiting-btc-research-refresh.service
[Unit]
Description=Leiting BTC research refresh job
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/leiting-btc
EnvironmentFile=-/opt/leiting-btc/btc_config/live.env
ExecStart=/opt/leiting-btc/scripts/btc/run_research_refresh.sh

[Install]
WantedBy=multi-user.target
