# deploy/btc/systemd/leiting-btc-engine.service
[Unit]
Description=Leiting BTC engine
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/leiting-btc
EnvironmentFile=-/opt/leiting-btc/btc_config/live.env
ExecStart=/opt/leiting-btc/scripts/btc/run_engine.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
