# deploy/btc/systemd/leiting-btc-review.timer
[Unit]
Description=Run Leiting BTC review periodically

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h
Unit=leiting-btc-review.service

[Install]
WantedBy=timers.target
