# deploy/btc/systemd/leiting-btc-demo-auto.service
[Unit]
Description=Leiting BTC demo auto execution service
After=network-online.target leiting-btc-engine.service leiting-btc-post-fill.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/leiting-btc
EnvironmentFile=-/opt/leiting-btc/btc_config/live.env
ExecStart=/opt/leiting-btc/scripts/btc/run_demo_auto_service.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
